# config/systemd/user/vnc_tunnel.service
[Unit]
Description=VNC Tunnel to Rocinante
Documentation=
After=network-online.target

[Service]
ExecStart=/usr/bin/ssh -o "ServerAliveInterval 30" -o "ServerAliveCountMax 3" -o ExitOnForwardFailure=yes -i /home/vanandrew/.ssh/id_rsa -v -N -L 3000:10.20.145.192:5908 vana@dosenbachlab.wustl.edu
Restart=always

[Install]
WantedBy=multi-user.target
